# systemd/rke2-airgap-doctor.service

[Unit]
Description=RKE2 Air-Gap Doctor (verify node health)
Documentation=file:/opt/rke2-airgap/docs/QUICKSTART.md

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'make -C /opt/rke2-airgap verify | tee -a /var/log/rke2-airgap-doctor.log'
Nice=10
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
